# bot/systemd/portfolio-bot.service
[Unit]
Description=Portfolio Telegram Bot
Documentation=https://github.com/cjmin-n/PORTFOLIO_REACT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# ─── 경로/유저 (Oracle Cloud의 'emma' 계정 기준) ─────────────────
# 다른 서버에서 쓸 경우 User와 모든 경로의 'emma' 부분을 해당 user로 교체
User=emma
WorkingDirectory=/home/emma/PORTFOLIO_REACT/bot
ExecStart=/home/emma/PORTFOLIO_REACT/bot/venv/bin/python bot.py

# ─── 자동 재시작 (다운 시 10초 후 재시작) ────────────────────────
Restart=always
RestartSec=10

# ─── 로그 ───────────────────────────────────────────────────
# journalctl -u portfolio-bot -f 로 보거나, 파일로 따로 저장하려면 아래 사용
# StandardOutput=append:/home/emma/PORTFOLIO_REACT/bot/bot.log
# StandardError=append:/home/emma/PORTFOLIO_REACT/bot/bot.log

[Install]
WantedBy=multi-user.target
